# packer/imms-build-agent/files/fix-agent-work-owner.service
[Unit]
Description=Fix agent work directory owner

[Service]
ExecStart=/etc/fix-agent-work-owner.sh

[Install]
WantedBy=multi-user.target